# webview/docker/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:apache2]
redirect_stderr=true
command=/usr/sbin/httpd -DFOREGROUND

[program:mysql]
redirect_stderr=true
command=/usr/bin/pidproxy /var/run/mysqld/mysqld.pid /usr/bin/mysqld_safe
