# script/copy_data.cron
*/1 * * * * /usr/local/cnipol_trunk/script/copy_data.sh
