# script/run_asym.cron
*/30 * * * * /eic/u/dsmirnov/cnipol/script/run_asym.sh
