# script/run_asym.cron
*/20 * * * * /usr/local/cnipol_trunk/script/run_asym.sh
