# script/runAsymR22.cron
*/30 * * * * /star/u/zchang/run22/cnipol/script/run_asym_r22.sh
