# script/runstar.cron
*/1 * * * * /star/u/zchang/run22/cnipol/script/runstar.csh >> /star/u/zchang/run22/cnipol/script/Dec03.log
